# setup/etc/systemd/system/apply.service
[Unit]
Description=Apply settings.yaml changes (scoped)

[Service]
Type=oneshot
ExecStart=/opt/rzans_vpn_main/settings/dispatcher.sh